# scripts/systemd/hashpool-prometheus-proxy.service
[Unit]
Description=Hashpool Prometheus (Translator Metrics)
After=network.target

[Service]
Type=simple
User=prometheus
Group=prometheus
WorkingDirectory=/opt/hashpool
ExecStartPre=+/bin/sh -c 'mkdir -p /var/lib/hashpool/prometheus-proxy && chown -R prometheus:prometheus /var/lib/hashpool/prometheus-proxy'
ExecStart=/usr/bin/prometheus \
  --config.file=/opt/hashpool/config/prometheus-proxy.yml \
  --storage.tsdb.path=/var/lib/hashpool/prometheus-proxy \
  --web.listen-address=127.0.0.1:9091
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
